# ansible/roles/ubuntu/files/home/.config/autostart/Conky.desktop
[Desktop Entry]
Exec=conky -d --display=:0 --quiet
Name=Conky
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
