# ansible/roles/ubuntu/files/home/.config/autostart/numlock ON.desktop
[Desktop Entry]
Exec=numlockx &
Name=numlock ON
Type=Application
Version=1.0
